# example.conf
[program:foo]
Name=foo
Path=C:\path\to\foo.exe
Exe=foo.exe
[program:bar]
Name=bar
Path=C:\path\to\bar.exe
Exe=bar.exe
[program:baz]
Name=baz
Path=C:\path\to\baz.exe
Exe=baz.exe
