# example.conf
[program:foo]
command=c:/example.exe
user=test
[program:bar]
command=c:/example2.exe
user=test
[program:baz]
command=c:/example3.exe
user=test2
